# adb/rootfs/etc/services.d/adb/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Android Debug Bridge
# Runs the Android Debug Bridge server
# ==============================================================================
declare -a options

bashio::log.info 'Starting the Android Debug Bridge server...'

# Start options for ADB
options+=(-a)
options+=(server)
options+=(nodaemon)

# Run ADB
exec adb "${options[@]}"
